# ulauncher/.config/autostart/ulauncher.desktop
[Desktop Entry]
Type=Application
Name=ULauncher
Exec=env GDK_BACKEND=x11 ulauncher --hide-window